# Cloudflare/crontab
#IPTABLES for cloudflare on startup and each month
@reboot /home/braam/./cloudflare_ips_cron.sh
0 0 1 * * /home/braam/./cloudflare_ips_cron.sh
